# DbAssistant/setup/systemd/dbtool-api.service
# systemd unit file for DbManagementTool REST API
#
# Install (Linux):
#   sudo cp systemd/dbtool-api.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable dbtool-api
#   sudo systemctl start  dbtool-api
#   journalctl -u dbtool-api -f

[Unit]
Description=DbManagementTool REST API
Documentation=https://github.com/dhananjay-chaturvedi/dbassistant
After=network.target dbtool-monitor.service

[Service]
Type=simple
User=%i

# ── adjust these paths ───────────────────────────────────────────────────────
WorkingDirectory=/opt/DbManagementTool
ExecStart=/opt/DbManagementTool/.venv/bin/python dbtool.py api \
    --host 0.0.0.0 \
    --port 8000
# ────────────────────────────────────────────────────────────────────────────

Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=dbtool-api

LimitNOFILE=65536
TimeoutStopSec=15

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/dbtool /tmp

[Install]
WantedBy=multi-user.target
